# crontab
0 6 * * * /Users/chrisjaimon/cron_scripts/spank_that_ass.sh
0 10 * * * /Users/chrisjaimon/cron_scripts/court_gila_he.sh
0 10 * * * /Users/chrisjaimon/cron_scripts/basketball_poll.sh
